# systemd/clawsquad-backup.timer
[Unit]
Description=Daily backup timer for ClawSquad database
Requires=clawsquad-backup.service

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=3600  # Random delay up to 1 hour

# Additional backup times (optional)
# OnCalendar=*-*-* 02:00:00
# OnCalendar=*-*-* 14:00:00

[Install]
WantedBy=timers.target